# mp_helperd.service
[Unit]
Description=MPD helper for batya's player
After=network.target

[Service]
Type=exec
